# rootfs/etc/services.d/wyoming_whisper_api_client/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the wyoming_whisper_api_client service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

DEBUG_LOGGING=$(bashio::config 'debug_logging')
API_URL=$(bashio::config 'api_url')/inference
PORT=10300
URI=tcp://0.0.0.0:$PORT

OPTIONS=""
if [ "$DEBUG_LOGGING" = true ]; then
  OPTIONS="$OPTIONS --debug"
fi

cd /wyoming-whisper-api-client
python3 -m wyoming_whisper_api_client --uri ${URI} --api ${API_URL} ${OPTIONS}
