# wyoming-whisper-api-client/rootfs/etc/s6-overlay/s6-rc.d/wyoming-whisper-api-client/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the wyoming-whisper-api-client
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

API_URL=$(bashio::config 'api_url')/inference
PORT=10300
URI=tcp://0.0.0.0:$PORT

FLAGS=""
if bashio::config.true 'debug_logging'; then
    FLAGS+=('--debug')
fi

bashio::log.info "Starting Wyoming Whisper API Client..."
cd /wyoming-whisper-api-client
exec python3 -m wyoming_whisper_api_client --uri ${URI} --api ${API_URL} ${FLAGS}
